# system/provisioning-server.service
[Unit]
Description=HA Panel Provisioning Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/ha-panel/repo/provisioning-ui
ExecStart=/opt/ha-panel/venv/bin/python server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
